# zfs-api.service
[Unit]
Description=ZFS API Server
After=network.target redis.service docker.service
Wants=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/zfs-api
Environment="PATH=/opt/zfs-api/venv/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/zfs-api/venv/bin/python /opt/zfs-api/zfs_api_server.py
Restart=always
RestartSec=10

# Security options
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zfs-api

[Install]
WantedBy=multi-user.target
